# conf/app.ini
[program: app]
command = python3 app.py
directory = /home/zucc_helper
autostart = true