# conf/beat.ini
[program: beat]
command = celery -A tasks beat -l info
directory = /home/zucc_helper
autostart = true