# conf/beat.ini
[program: beat]
command = celery -A tasks beat
directory = /home/zucc_helper
autostart = true