# conf/celery.ini
[program: celery]
command = celery -A tasks worker -c 8
directory = /home/zucc_helper
autostart = true