# conf/celery.ini
[program: celery]
command = celery -A tasks worker -l info -c 8
directory = /home/zucc_helper
autostart = true